# systemd/curytybainboxd.service
[Unit]
Description=curytybainboxd
After=network.target

[Service]
EnvironmentFile=/etc/curytybainbox/environment
Restart=always
ExecStart=/usr/bin/curytybainboxd

[Install]
WantedBy=multi-user.target
